# contrib/kaijutsu-server.service
[Unit]
Description=Kaijutsu Server (kaish execution + collaboration)
Documentation=https://github.com/tobert/kaijutsu
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash -c 'export ANTHROPIC_API_KEY=$(< %h/.anthropic-key.txt) && exec %h/src/kaijutsu/target/debug/kaijutsu-server'
WorkingDirectory=%h/src/kaijutsu
Restart=on-failure
RestartSec=5

# Environment
Environment=RUST_LOG=info
Environment=OTEL_EXPORTER_OTLP_ENDPOINT=http://localhost:4317
Environment=OTEL_EXPORTER_OTLP_PROTOCOL=grpc
Environment=OTEL_SERVICE_NAME=kaijutsu-server

[Install]
WantedBy=default.target
