# contrib/kaijutsu-server.service
[Unit]
Description=Kaijutsu Server (kaish execution + collaboration)
Documentation=https://github.com/tobert/kaijutsu
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash -c 'export ANTHROPIC_API_KEY=$(cat %h/.anthropic-key.txt) && exec cargo run --release -p kaijutsu-server --manifest-path %h/src/kaijutsu/Cargo.toml'
WorkingDirectory=%h/src/kaijutsu
Restart=on-failure
RestartSec=5

# Environment
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
